# zookeeper/zookeeper-supervisord.conf
[program:zookeeper]
command=/usr/share/zookeeper/bin/zkServer.sh start

