# zookeeper/zookeeper-supervisord.conf
[program:zookeeper]
command=bin/zkServer.sh start

